# roles/config/templates/rc.local
/usr/sbin/nginx
/usr/local/php/sbin/php-fpm

exit 0